# platform/supervisord.conf
[supervisord]
nodaemon=true

[program:cassandra]
command=/cassandra-configurator.sh cassandra -f

[program:sparkWorker]
command=/start-worker.sh
